# etc/cron.hourly/bedtime
#!/bin/bash

/usr/share/bedtime/bin/update_ipt
/usr/share/bedtime/bin/dhcp_slurp
