# usr/share/bedtime/bedtime.service
[Unit]
Description=Bedtime

[Service]
Type=forking
ExecStart=/usr/share/bedtime/bin/bedtimed

[Install]
WantedBy=multi-user.target
